# contrib/go-glpi-agent.service
[Unit]
Description=Go GLPI Agent (inventory run)
Documentation=https://github.com/jniltinho/go-fusioninventory-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# `run` does a single inventory cycle (server or local, per agent.cfg).
# Scheduling is handled by the companion .timer.
ExecStart=/opt/go-glpi-agent/go-glpi-agent run --conf-file /opt/go-glpi-agent/agent.cfg
# runs as root to read dmidecode, /sys/class/dmi and disk serials
User=root

[Install]
WantedBy=multi-user.target
