# contrib/macos/com.glpi.go-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.glpi.go-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/go-glpi-agent/go-glpi-agent</string>
        <string>run</string>
        <string>--conf-file</string>
        <string>/usr/local/go-glpi-agent/agent.cfg</string>
    </array>
    <!-- Run once at load, then every 3600s (the macOS analog of the systemd timer). -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/var/log/go-glpi-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/go-glpi-agent.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
